# ita/roles/stornext_fc_inst/templates/tsih2/crontab
@reboot /site/maint/bin/nfs_restart.csh > /site/maint/var/nfs_restart.log 2>&1
